# systemd/bing-wallpaper.timer
TIMER_FILE="$HOME/.config/systemd/user/bing-wallpaper.timer"

cat > "$TIMER_FILE" <<EOF
[Unit]
Description=Bing wallpaper timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=$INTERVAL
Persistent=true

[Install]
WantedBy=timers.target
EOF