# systemd/bing-wallpaper.timer
[Unit]
Description=Run Bing Wallpaper every 6 hours

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target